# packer/files/crontab
SHELL=/bin/bash
MAILTO=""

05 06 * * * /root/jaaxman/scripts/fetchrss.sh >> /var/log/jaaxman/crontab.log
