# packer/files/crontab
SHELL=/bin/bash
PATH=/root/local/python-3.6.2/bin:/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

05 06 * * * python /root/jaaxman/manage.py fetchrss >> /var/log/jaaxman/crontab.log
